# files/etc/uci-defaults/99-custom

# Enable DHCP client on lan
uci set network.lan.proto="dhcp"

# Enable WAN over serial port for SIM7070G
uci set network.wwan0=interface
uci set network.wwan0.device='/dev/ttyAMA0'
uci set network.wwan0.apn='internet.m2mportal.de'
uci set network.wwan0.service='umts'
uci set network.wwan0.proto='3g'

# Optional WAN for a USB CDC modem
uci set network.wwan1=interface
uci set network.wwan1.proto='dhcp'
uci set network.wwan1.ifname='usb0'

uci commit network
service network restart

uci add_list firewall.@zone[1].network='wwan0'
uci add_list firewall.@zone[1].network='wwan1'

uci commit firewall
service firewall restart
